# sv/brightnessctld/run
#!/bin/sh
[ -r brillo ] && brightnessctl -q set "$(cat brillo)"
exec chpst -b brightnessctl pause
